# launchd/com.yssemanticsystems.powera-dolphin-pipe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yssemanticsystems.powera-dolphin-pipe</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/powera-dolphin-pipe</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/var/log/powera-dolphin-pipe.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/powera-dolphin-pipe.log</string>
</dict>
</plist>

